# docs/examples/nautica-api.service
[Unit]
Description=Nautica API
After=network.target

[Service]
Type=simple
User=nautica
Group=nautica
WorkingDirectory=/websites/nautica/backend
EnvironmentFile=/etc/nautica/nautica.env
ExecStart=/usr/bin/node /websites/nautica/backend/server-dist/apps/api/src/server.js
Restart=always
RestartSec=5
TimeoutStartSec=120
UMask=0027
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/websites/nautica/shared/db /websites/nautica/shared/uploads
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
